# autostart/chrome.desktop
[Desktop Entry]
Type=Application
Exec=google-chrome
