# autostart/radiotray.desktop
[Desktop Entry]
Type=Application
Exec=radiotray
