# autostart/terminal.desktop
[Desktop Entry]
Type=Application
Exec=gnome-terminal
